# APP_Source/Network/tailscale/input/all/etc/systemd/system/tailscaled.service
[Unit]
Description=Tailscale node agent
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/bin/bash /etc/tailscale/tailscale-init.sh
Restart=always
RestartSec=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target